# sample.ini
[program:hello]
directory=/home/ec2-user/WebRTCHandsOnSig
command=go run *.go
user=root
autorestart=true
stdout_logfile=/home/ec2-user/WebRTCHandsOnSig/hello.log
redirect_stderr=true
